# deploy/rn1-paper-surfandturf.service
[Unit]
Description=Surfandturf Paper — copy of wallet 0x9f2fe025…2ca8 NBA/foot conviction bot
After=network.target bot-cp.service

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/bot-trading
EnvironmentFile=-/home/botuser/bot-trading/.env
# Surf is a low-volume, high-conviction whale (avg chunk $5K-$150K). Backtest
# 2026-05-20 on his 9 sold-and-tracked positions (48d window) showed 66.7% WR
# at +64.9% ROI per trade. At $4.50 sizing the equity barely moves (his ~1 trade
# /week × $0.60 EV = invisible). Scaled to $20 normal / $4 penny + max 5 pos
# to make each copy meaningful while keeping max-blast acceptable (-$100 = 17%).
Environment=SURF_INITIAL=600
Environment=SURF_MAX_POSITIONS=5
Environment=SURF_MIN_TARGET=500
Environment=SURF_MAX_USD_PER_MARKET=20
Environment=SURF_SKIP_HIGH=0.20
Environment=SURF_NORMAL_SIZE=20
Environment=SURF_PENNY_SIZE=4
Environment=SURF_NORMAL_MAX=0.85
ExecStart=/home/botuser/bot-trading/venv/bin/python -m analysis.rn1.paper_surfandturf
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/botuser/bot-trading/logs/rn1/paper_surfandturf.log
StandardError=append:/home/botuser/bot-trading/logs/rn1/paper_surfandturf.log

[Install]
WantedBy=multi-user.target
