# deploy/rn1-paper-absband.service
[Unit]
Description=RN1 Paper Absband — "yesterday strategy" small sizing, no Option B filters
After=network.target bot-cp.service

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/bot-trading
EnvironmentFile=-/home/botuser/bot-trading/.env
# "Yesterday's strategy" — small sizing that grew $30 → $60 on surfandturf,
# now applied to RN1 with $600 paper bankroll. No Option B filters; raw absband.
Environment=RN1_ABSBAND_INITIAL=600
Environment=RN1_ABSBAND_MAX_POSITIONS=8
Environment=RN1_ABSBAND_MIN_TARGET=200
Environment=RN1_ABSBAND_MAX_USD_PER_MARKET=5
Environment=RN1_ABSBAND_SKIP_HIGH=0.20
Environment=RN1_ABSBAND_NORMAL_SIZE=4.5
Environment=RN1_ABSBAND_PENNY_SIZE=1.0
ExecStart=/home/botuser/bot-trading/venv/bin/python -m analysis.rn1.paper_absband
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/botuser/bot-trading/logs/rn1/paper_absband.log
StandardError=append:/home/botuser/bot-trading/logs/rn1/paper_absband.log

[Install]
WantedBy=multi-user.target
